# samples/init/launchd/virtualbox.jenkins.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>virtual.jenkins</string>

	<key>ProgramArguments</key>
	<array>
		<string>VBoxManage</string>
		<string>startvm</string>
		<string>Tout Jenkins Server</string>
		<string>-type</string>
		<string>headless</string>
	</array>

	<key>UserName</key>
	<string>tout</string>

	<key>WorkingDirectory</key>
	<string>/Users/tout/VirtualBox</string>

	<key>RunAtLoad</key>
	<true/>

	<key>StartInterval</key>
	<integer>60</integer>

	<key>StandardErrorPath</key>
	<string>/Users/tout/VirtualBox</string>

	<key>StandardOutPath</key>
	<string>/Users/tout/VirtualBox</string>
</dict>
</plist>
